# root/etc/profile.d/one_time_setup_disabled.sh
#!/bin/sh

mv $(dirname $0)/one_time_setup.sh $(dirname $0)/one_time_setup.sh.disabled
